# ente.service
[Service]
Type=simple
ExecStart=/usr/bin/python /root/pidora-keyboard.py
